# home-justin.wilson-nfs.automount
[Automount]
Where=/home/justin.wilson/nfs

[Install]
WantedBy=graphical.target
