# home-justin.wilson-Dropbox.automount
[Automount]
Where=/home/justin.wilson/Dropbox

[Install]
WantedBy=graphical.target
